# root/laptop/etc/udev/rules.d/90-backlight.rules
SUBSYSTEM=="backlight", ACTION=="add", \
	RUN+="/bin/chgrp users %S%p/brightness", \
	RUN+="/bin/chmod g+w %S%p/brightness"
