# files/gui/.config/systemd/user/konawall.timer
[Timer]
OnBootSec=10
OnUnitInactiveSec=20min

[Install]
WantedBy=timers.target
