# files/personal/.config/systemd/user/tasksync.service
[Service]
Type=simple
ExecStart=/usr/bin/env task sync
