# root/shanghai/etc/systemd/system/flandre-ssh.socket
[Socket]
ListenStream=33022
Accept=true

[Install]
WantedBy=sockets.target
