# root/shanghai/etc/systemd/system/remote-backup.timer
[Timer]
OnCalendar=*-*-* 20:30:00
Persistent=true

[Install]
WantedBy=timers.target
