# root/shanghai/etc/systemd/system/vpsdump.path
[Path]
DirectoryNotEmpty=/run/vpsdump
MakeDirectory=true
DirectoryMode=0777
