# files/shanghai/.config/systemd/user/getquote.service
[Service]
Type=simple
ExecStart=/usr/bin/env ${HOME}/projects/gensokyo/ledger/update_prices
