# root/base/etc/systemd/system/hdd-apm@.service
[Unit]
Description=Hard drive power saving

[Service]
Type=oneshot
ExecStart=/sbin/hdparm -B 251 /dev/disk/by-id/wwn-%i

[Install]
WantedBy=multi-user.target
