# root/base/etc/systemd/system/inactive-shutdown.timer
[Timer]
OnBootSec=30m
OnUnitInactiveSec=10m

[Install]
WantedBy=timers.target
