# root/base/etc/systemd/system/sshd-local-proxy@.socket
[Socket]
ListenStream=%i
Service=sshd-local-proxy.service

[Install]
WantedBy=sockets.target
