# root/hw-x370gpc/etc/systemd/system/cpuclock.service
[Service]
Type=oneshot
ExecStart=/usr/bin/zenstates -p 0 -f 9e -d 8 -v 1e
