# files/base/.config/systemd/user/dotfiles.service
[Service]
Type=simple
ExecStart=/usr/bin/env ${HOME}/.dotfiles/dot.sh update-self
