# root/satorin/etc/systemd/system/connman.service.d/disable_if.conf
[Service]
ExecStart=
ExecStart=/usr/bin/connmand -n --nodnsproxy -I eth0 -I wlan0
