# root/satorin/etc/systemd/system/ipsecmon.service
[Unit]
After=connman.service

[Service]
Type=simple
Restart=always
ExecStart=/root/.bin/_/ipsecmon
